# etc/cron.d/update-time
*/5 * * * * root ntpdate us.pool.ntp.org >/dev/null 2>&1
